# roles/systemd-timers/templates/fakedhcp.timer.j2
[Unit]
Description=Fake DHCP client for MTS ISP

[Timer]
OnBootSec=0s
OnUnitActiveSec=5m 

[Install]
WantedBy=timers.target
