# roles/systemd-timers/templates/ping@.timer.j2
[Unit]
Description=ping timer
PartOf=keepalived-gateway.target

[Timer]
OnBootSec=0s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
