# roles/keepalived/templates/systemd/keepalived-gateway.target.j2
[Unit]
Description=keepalived controlled target
Requires=fakedhcpc@isp2.service
