# roles/keepalived/templates/systemd/ping@.timer.j2
[Unit]
Description=ping timer
PartOf=keepalived-gateway.target

[Timer]
OnActiveSec=0s
OnUnitActiveSec=15s
AccuracySec=1s

[Install]
WantedBy=timers.target
